# stratosphere-dist/src/deb/bin/jobmanager
#! /bin/sh
########################################################################################################################
# Copyright (C) 2010-2013 by the Stratosphere project (http://stratosphere.eu)
#
# Licensed under the Apache License, Version 2.0 (the "License"); you may not use this file except in compliance with
# the License. You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software distributed under the License is distributed on
# an "AS IS" BASIS, WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied. See the License for the
# specific language governing permissions and limitations under the License.
########################################################################################################################

#
# skeleton  example file to build /etc/init.d/ scripts.
#    This file should be used to construct scripts for /etc/init.d.
#
#    Written by Miquel van Smoorenburg <miquels@cistron.nl>.
#    Modified for Debian
#    by Ian Murdock <imurdock@gnu.ai.mit.edu>.
#               Further changes by Javier Fernandez-Sanguino <jfs@debian.org>
#
# Version:  @(#)skeleton  1.9  26-Feb-2001  miquels@cistron.nl
#
# Starts a stratosphere jobmanager
#
# chkconfig: 2345 85 15
# description: stratosphere jobmanager
#
### BEGIN INIT INFO
# Provides:          stratosphere-jobmanager
# Required-Start:    $network $local_fs
# Required-Stop:
# Should-Start:      $named
# Should-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: stratosphere jobmanager daemon
### END INIT INFO

# Include stratosphere defaults if available
if [ -f /etc/default/stratosphere ] ; then
  . /etc/default/stratosphere
fi


if [ "$STRATOSPHERE_PID_DIR" = "" ]; then
	STRATOSPHERE_PID_DIR=/tmp
fi

if [ "$STRATOSPHERE_IDENT_STRING" = "" ]; then
	STRATOSPHERE_IDENT_STRING="$USER"
fi

STRATOSPHERE_HOME=/usr/share/stratosphere-dist
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DAEMON_SCRIPT=$STRATOSPHERE_HOME/bin/jobmanager.sh
NAME=stratosphere-jobmanager
DESC="stratosphere jobmanager daemon"
PID_FILE=$STRATOSPHERE_PID_DIR/stratosphere-$STRATOSPHERE_IDENT_STRING-jobmanager.pid

test -x $DAEMON_SCRIPT || exit 1


DODTIME=5                   # Time to wait for the server to die, in seconds
                            # If this value is set too low you might not
                            # let some servers to die gracefully and
                            # 'restart' will not work

# Checks if the given pid represents a live process.
# Returns 0 if the pid is a live process, 1 otherwise
stratosphere_is_process_alive() {
  local pid="$1" 
  ps -fp $pid | grep $pid | grep jobmanager > /dev/null 2>&1
}

# Check if the process associated to a pidfile is running.
# Return 0 if the pidfile exists and the process is running, 1 otherwise
stratosphere_check_pidfile() {
  local pidfile="$1" # IN
  local pid

  pid=`cat "$pidfile" 2>/dev/null`
  if [ "$pid" = '' ]; then
    # The file probably does not exist or is empty. 
    return 1
  fi
  
  set -- $pid
  pid="$1"

  stratosphere_is_process_alive $pid
}

stratosphere_process_kill() {
   local pid="$1"    # IN
   local signal="$2" # IN
   local second

   kill -$signal $pid 2>/dev/null

   # Wait a bit to see if the dirty job has really been done
   for second in 0 1 2 3 4 5 6 7 8 9 10; do
      if stratosphere_is_process_alive "$pid"; then
         # Success
         return 0
      fi

      sleep 1
   done

   # Timeout
   return 1
}

# Kill the process associated to a pidfile
stratosphere_stop_pidfile() {
   local pidfile="$1" # IN
   local pid

   pid=`cat "$pidfile" 2>/dev/null`
   if [ "$pid" = '' ]; then
      # The file probably does not exist or is empty. Success
      return 0
   fi
   
   set -- $pid
   pid="$1"

   # First try the easy way
   if stratosphere_process_kill "$pid" 15; then
      return 0
   fi

   # Otherwise try the hard way
   if stratosphere_process_kill "$pid" 9; then
      return 0
   fi

   return 1
}

start() {
    $STRATOSPHERE_HOME/bin/jobmanager.sh start cluster
}

stop() {
    $STRATOSPHERE_HOME/bin/jobmanager.sh stop
}

check_for_root() {
  if [ $(id -ur) -ne 0 ]; then
    echo 'Error: root user required'
    echo
    exit 1
  fi
}

stratosphere_service() {
    case "$1" in
         start)
            check_for_root
            echo -n "Starting $DESC: "
            start

            if stratosphere_check_pidfile $PID_FILE ; then
                echo "$NAME."
            else
                echo "ERROR. Could not start $DESC"
                exit 1
            fi
            ;;
        stop)
            check_for_root
            echo -n "Stopping $DESC: "
            stop
            [ -n "$DODTIME" ] && sleep $DODTIME

            if stratosphere_check_pidfile $PID_FILE ; then
                echo "ERROR. Could not stop $DESC"
                exit 1
            else
                echo "$NAME."
            fi
            ;;
        force-stop)
            check_for_root
            echo -n "Forcefully stopping $DESC: "
            stratosphere_stop_pidfile $PID_FILE
            [ -n "$DODTIME" ] && sleep $DODTIME

            if ! stratosphere_check_pidfile $PID_FILE ; then
                echo "$NAME."
            else
                echo "ERROR. Could not force stop $DESC"
                exit 1
            fi
            ;;
        force-reload)
            check_for_root
            echo -n "Forcefully reloading $DESC: "
            stratosphere_check_pidfile $PID_FILE && $0 restart
            ;;
        restart)
            check_for_root
            echo -n "Restarting $DESC: "
            stop
            [ -n "$DODTIME" ] && sleep $DODTIME
            $0 start
            ;;
        status)
            echo -n "$NAME is "
            if stratosphere_check_pidfile $PID_FILE ;  then
                echo "running"
            else
                echo "not running."
                exit 1
            fi
            ;;
        *)
            N=/etc/init.d/$NAME
            echo "Usage: $N {start|stop|restart|force-reload|status|force-stop}" >&2
            exit 1
            ;;
    esac
}

stratosphere_service "$1"

exit 0
